# playbooks/roles/slurm/templates/systemd/munge.service.d/unit.conf.j2
[Service]
ExecStart=
ExecStart=/usr/sbin/munged  --key-file {{ munge_conf_path }}/munge.key

