# playbooks/roles/weka_client/templates/weka-snapshot-umount.timer.j2
[Unit]
Description=Daily Snapshot unmounting
Requires=weka-snapshot-umount.service

[Timer]
OnCalendar=*-*-* 7:00:00

[Install]
WantedBy=timers.target
