# playbooks/roles/weka_client/templates/weka-snapshot-automounting.timer.j2
[Unit]
Description=Daily Snapshot Automounting
Requires=weka-snapshot-automounting.service

[Timer]
OnCalendar=*-*-* 13:00:00

[Install]
WantedBy=timers.target
