# playbooks/roles/slurm/templates/systemd/slurmctld.service.d/unit.conf.j2
[Unit]
Requires=munge.service {{ "slurmdbd.service" if not slurm_federation }}

[Service]
Environment="SLURMCTLD_OPTIONS=-R"
Restart=always
