# playbooks/roles/slurm/templates/systemd/slurmctld.service.d/unit.conf.j2
[Unit]
Requires=munge.service slurmdbd.service

[Service]
Environment="SLURMCTLD_OPTIONS=-R"
Restart=always
