# playbooks/roles/mongoDB/templates/systemd/node_updater.timer
[Unit]
Description=Run Node discovery every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target