# playbooks/roles/slurm_rest_api/templates/slurm-exporter.service.j2
[Unit]
Description=Slurm Exporter
Wants=network-online.target
After=network-online.target

[Service]
User={{ slurm_user }}
Group={{ slurm_group }}
Type=simple

# Activate the virtual environment before running
Environment="PATH={{ venv_dir }}/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONPATH={{ venv_dir }}/lib/python{{ python_version }}/site-packages"

# Run exporter from virtual environment
ExecStart={{ venv_dir }}/bin/python {{ venv_dir }}/bin/slurm-exporter.py

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target