# playbooks/roles/slurm/templates/systemd/slurmdbd.service.d/unit.conf.j2
[Unit]
Requires=munge.service

[Service]
Restart=always
