# playbooks/roles/metrics-exporter/templates/dcgm-exporter.service.j2
[Unit]
Description=NVIDIA DCGM Exporter
After=network.target

[Service]
ExecStart=/usr/bin/dcgm-exporter -f /etc/dcgm-counters.csv

[Install]
WantedBy=multi-user.target
