# playbooks/roles/nodes_update/templates/systemd/node_updater.timer
[Unit]
Description=Run Node discovery every 2 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target